# scripts/update_dashboard.crontab
# every Monday, fetch data and update html in noaa-onms/climate-dashboard
0 0 * * 1 /share/github/noaa-onms/climate-dashboard/scripts/update_dashboard.sh >> /var/log/cron.log 2>&1
